# octolab/octolab.service
# Process-level supervision for octolab (the "via systemd" half).
#
# Install (user service):
#   cp octolab.service ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now octolab
#   journalctl --user -u octolab -f
#
# `Restart=always` resurrects the process on ANY exit, including a clean one —
# so the agent can *restart itself after applying config* simply by exiting
# (the gap OpenClaw has). In-process supervision (octo-core) restarts a single
# flaky connector without a full restart; this handles the whole-process case.

[Unit]
Description=octolab — Octo runtime agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/nikita/code/personal/octo/octo/octolab
# The binary loads the repo-root .env itself; or point systemd at it directly:
# EnvironmentFile=/home/nikita/code/personal/octo/.env
# Binary lives in the workspace target dir (octo/target), not the crate dir.
ExecStart=/home/nikita/code/personal/octo/octo/target/release/octolab
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
